# src/config/at-base.service.in
[Unit]
Description=Away-Team VPN Server

[Service]
ExecStart=/usr/sbin/AT/at-base.py

[Install]
WantedBy=multi-user.target
